# examples/whisper-pod-package/deploy/whisper.container
[Container]
Image=docker.io/library/asmo-whisper:latest
PublishPort=8001:8001
