# deploy/systemd/asmo-pod-orchestrator.service
[Unit]
Description=Asmo Pod Orchestrator
After=default.target

[Service]
Type=simple
WorkingDirectory=/home/clohl/src/asmo-pod-orchestrator
Environment=PATH=/home/clohl/.nvm/versions/node/v24.14.0/bin:/usr/local/bin:/usr/bin:/bin
Environment=HOST=0.0.0.0
Environment=PORT=8787
ExecStart=/home/clohl/.nvm/versions/node/v24.14.0/bin/npm start
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
